# files/etc/uci-defaults/89_lime-release
#!/bin/sh

# override original lime release because it overwrites buildroot release information.

exit 0
